# supervisord.conf
[program:reactivebb]
command=python3.11 /home/rulo/reactive-bb-build/start.py
autorestart = true

[supervisord]
